# deploy/mdm-purge.service
[Unit]
Description=MDM retention purge (one-shot)

[Service]
Type=oneshot
User=mdm-svc
Group=mdm-svc
WorkingDirectory=__REPO_DIR__
ExecStart=__REPO_DIR__/.venv/bin/python -m mdm.retention_purge

# Hardening — same posture as mdm.service
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ReadWritePaths=__REPO_DIR__/data
